# finance-frontend.service
[Unit]
Description=Finance Planner Frontend (npm)
After=network.target
Wants=network.target

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/var/www/finance
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 3000
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=finance-frontend

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/finance

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target 